# overlays/tl-wr740n-v4-vpn-wds/files/etc/profile.d/tp-link-fw.sh
export TP_LINK_FW_BUILD="wr740n-v4-vpn-wds"
